# upstart/sync-monitor.conf.in
description "sync-monitor"
author "Renato Oliveira Filho <renato.filho@canonical.com>"

start on started address-book-service
stop on session-end

respawn

exec @CMAKE_INSTALL_FULL_LIBEXECDIR@/sync-monitor
